# deploy/ai.caihub.ai-food-news-weekly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.caihub.ai-food-news-weekly</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>/Users/ming/CaiHub/scripts/run_ai_food_news_weekly_email.sh</string>
  </array>

  <key>RunAtLoad</key>
  <false/>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key>
    <integer>1</integer>
    <key>Hour</key>
    <integer>10</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/ming/.openclaw/logs/ai-food-news-weekly.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/ming/.openclaw/logs/ai-food-news-weekly.err.log</string>
</dict>
</plist>
